# supervisor/chromium.conf
[program:chromium]
priority=100
user=%USER%
command=/start-chromium.sh
environment=DISPLAY=":1"
autorestart=true