# supervisor/chromium.conf
[program:chromium]
priority=100
user=%USER%
command=chromium-browser --remote-debugging-port=7200 --disable-gpu --no-sandbox --start-maximized
environment=DISPLAY=":1"
autorestart=true