# supervisor/socat-chromium.conf
[program:socat]
command=/usr/bin/socat tcp-listen:7201,fork tcp:localhost:7200